# ansible/roles/ebpf_exporter/files/ebpf-execsnoop.service
[Unit]
Description=eBPF execsnoop tracer (bpftrace -> log file for Loki/Promtail)
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/stdbuf -oL /usr/bin/bpftrace /opt/ebpf-snoop-tools/execsnoop.bt
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/ebpf/execsnoop.log
StandardError=journal

[Install]
WantedBy=multi-user.target